# click/deployment/systemd/gate.service
[Unit]
Description=JIVO Gate Management (Django)
After=network.target

[Service]
User=superadmin
Group=www-data
WorkingDirectory=/home/superadmin/django_projects/Task/click
# Served under /gate/; nginx strips the prefix, SCRIPT_PREFIX restores it for URL generation.
Environment=SCRIPT_PREFIX=/gate
# Reuse the Team Board project .env (DB creds, ALLOWED_HOSTS). settings.py also
# loads ../.env directly, so this is optional (leading '-').
EnvironmentFile=-/home/superadmin/django_projects/Task/.env
ExecStart=/home/superadmin/django_projects/Task/click/.venv/bin/gunicorn config.wsgi:application --bind 127.0.0.1:8731 --workers 3 --timeout 60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
