# click/deployment/systemd/gate.service
[Unit]
Description=JIVO Gate Management (Django)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/teamboard/click
EnvironmentFile=/var/www/teamboard/click/.env
ExecStart=/var/www/teamboard/click/.venv/bin/gunicorn config.wsgi:application --bind 127.0.0.1:8731 --workers 3 --timeout 60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
